# systemd/koast-reviews-sync.service
[Unit]
Description=Koast Reviews Sync Worker
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/koast-workers
ExecStart=/usr/bin/python3 /home/ubuntu/koast-workers/reviews_sync.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
