# systemd/koast-pricing-worker.service
[Unit]
Description=Koast Pricing Worker
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/koast-workers
ExecStart=/usr/bin/python3 /home/ubuntu/koast-workers/pricing_worker.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
