# systemd/koast-pricing-validator.timer
[Unit]
Description=Run Koast pricing validator daily at 10:00 UTC

[Timer]
OnCalendar=*-*-* 10:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
