# systemd/koast-messaging-executor.service
[Unit]
Description=Koast messaging template executor — draft-and-queue (Session 8a)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/staycommand-workers
EnvironmentFile=/home/ubuntu/staycommand-workers/.env
ExecStart=/home/ubuntu/staycommand-workers/.venv/bin/python /home/ubuntu/staycommand-workers/messaging_executor.py
StandardOutput=append:/var/log/staycommand/messaging-executor.log
StandardError=append:/var/log/staycommand/messaging-executor.log

[Install]
WantedBy=multi-user.target
